# ansible/vultr/cointhink-api.service.jinja
[Unit]
Description=Cointhink API

[Service]
ExecStart=/home/devops/api-go/cointhink
WorkingDirectory=/home/devops/api-go
